# baseline-el9/files/el9/etc/profile.d/halt.sh
#Owner:           root
#Group:           root
#Mode:            644
alias halt='echo "Do not do that here"' 2>/dev/nul

